# service/coexist-weekly.timer
[Unit]
Description=weekly coexist backup script

[Timer]
OnCalendar=weekly
AccuracySec=1us
Persistent=true

[Install]
WantedBy=default.target
