# service/coexist-daily.timer
[Unit]
Description=daily coexist backup script

[Timer]
OnCalendar=daily
AccuracySec=1us
Persistent=true

[Install]
WantedBy=default.target
